# deploy/gcp-vm1/systemd/smc-demo-runner.service
[Unit]
Description=SMC OB+FVG strategy demo runner
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=aungp
WorkingDirectory=/opt/session-smc-trading-bot
EnvironmentFile=/etc/session-smc-trading-bot/live-dashboard.env
ExecStart=/opt/session-smc-trading-bot/deploy/gcp-vm1/run_smc_demo.sh
Restart=always
RestartSec=15
Nice=5

[Install]
WantedBy=multi-user.target
