# deploy/gcp-vm1/systemd/vps-health-check.timer
[Unit]
Description=Run VPS health check every 15 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=15min
Persistent=true
Unit=vps-health-check.service

[Install]
WantedBy=timers.target
